# wazuh/config/etc/service/api/run
#!/usr/bin/bash
echo >&2 "starting API"

sleep 5
echo "START API"
exec /bin/node /var/ossec/api/app.js
