# wazuh/config/etc/services.d/manager/run
#!/usr/bin/with-contenv bash
echo >&2 "starting Manager"

/var/ossec/bin/ossec-control start
sleep infinity
