# wazuh/services/startup/finish
#!/usr/bin/env sh
echo >&2 "service1 exited. code=${1}"
