# wazuh/services/api/run
#!/usr/bin/with-contenv sh
echo >&2 "starting API"

sleep 5
/bin/node /var/ossec/api/app.js
