# wazuh/services/manager/finish
#!/usr/bin/env bash
echo >&2 "Manager exited. code=${1}"

# terminate other services to exit from the container
#exec s6-svscanctl -t /var/run/s6/services

